# reference/zapret2/init.d/runit/zapret2/finish
#!/bin/sh
/opt/zapret2/init.d/sysv/zapret2 stop
