# reference/zapret2/init.d/runit/zapret2/run
#!/bin/sh
/opt/zapret2/init.d/sysv/zapret2 start
exec chpst -b zapret2 sleep infinity
